# docker/data/crontabs/app
 * * * * * cd /apps/APP_NAME/ && ./yii cron/run
